# contrib/multi-window-bridge/launchd/com.local.mwb-watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- __HOME__ is substituted with $HOME by mwb-install at install time. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.local.mwb-watchdog</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>__HOME__/.claude/multi-window-bridge/daemon/watchdog.py</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>__HOME__</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>ProcessType</key>
    <string>Background</string>

    <key>StandardOutPath</key>
    <string>__HOME__/.claude/multi-window-bridge/logs/watchdog.out.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/.claude/multi-window-bridge/logs/watchdog.err.log</string>
</dict>
</plist>
